# supervisor/flask-pypi-proxy.conf
[program:flask-pypi-proxy]
user=www-data
command=/home/ubuntu/.virtualenvs/flask-pypi-proxy/bin/uwsgi --ini /data/web/flask-pypi-proxy/uwsgi.ini
autostart=true
